# supervisor-app.conf.tmpl
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini {{ .Env.UWSGI_INI_PATH }}

[program:nginx-PictureShareService]
command = /usr/sbin/nginx